# install/roles/elk-metricbeat.arm/templates/metricbeat.service
[Unit]
Description=metricbeat
Documentation=https://www.elastic.co/guide/en/beats/metricbeat/current/index.html
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/opt/metricbeat/metricbeat-linux-arm -c /etc/metricbeat/metricbeat.yml
Restart=always

[Install]
WantedBy=multi-user.target